# tatooine/etc/systemd/system/clean-conda-cache.timer
[Unit]
Description=Check and clean conda pkgs cache

[Timer]
OnActiveSec=2d
Unit=clean-conda-cache.service

[Install]
WantedBy=multi-user.target
